# ha_addon/evcharging_tracker/rootfs/etc/services.d/evcharging_tracker/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting EV Charging Tracker..."

cd /app
python run_proxy.py